# systemd/plasma-agent.service
[Unit]
Description=plasma-agent (backup service)

[Service]
Type=simple
ExecStart=/usr/bin/plasma-agent
StandardError=journal
